# deploy/ama-listing.service
[Unit]
Description=AMA Listing Creator (Next.js)
After=network.target

[Service]
Type=simple
User=portal-system
WorkingDirectory=/opt/portal-system/apps/ama-listing-creator/.next/standalone
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

Environment=NODE_ENV=production
Environment=PORT=3002
Environment=HOSTNAME=127.0.0.1
# Persistent data dir lives at project root (outside the standalone bundle that
# gets rebuilt on every deploy). listingStore.getDataRoot() reads this.
Environment=AMA_LISTING_DATA_ROOT=/opt/portal-system/apps/ama-listing-creator/data

[Install]
WantedBy=multi-user.target
